# deploy/frp/frpc.service.example
[Unit]
Description=frp client for Wan inference API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/frp
ExecStart=/opt/frp/frpc -c /opt/frp/frpc.toml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
